# conf/dynette.cron
* * * * * #USER# /usr/bin/python #DESTDIR#/dynette.cron.py >> /var/log/dynette.log 2>&1
