# conf/dynette-regen-named-conf.path
[Path]
Unit=__APP__-regen-named-conf.service
PathChanged=__DATA_DIR__

[Install]
WantedBy=multi-user.target
